# scripts/ops/revampit-cron@close-decisions.timer
[Unit]
Description=Schedule RevampIT cron: close-decisions (daily 00:00 UTC)

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
